# systemd/autoclaude.timer
[Unit]
Description=AutoClaude Poller — runs every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=autoclaude.service

[Install]
WantedBy=timers.target
